# configs/systemd/user/notmuch.service
[Unit]
Description=Notmuch index service

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new --quiet
