# configs/systemd/user/notmuch.timer
[Unit]
Description=Notmuch index timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=10s
Unit=notmuch.service

[Install]
WantedBy=timers.target
